# ham-cw.service
[Unit]
Description=ham-cw iambic CW keyer
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/Ham-CW
ExecStart=/usr/bin/python3 /home/pi/Ham-CW/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
