# ham-cw.service
[Unit]
Description=ham-cw CW keyer
After=network.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/home/pi/Ham-CW
ExecStart=/usr/bin/python3 /home/pi/Ham-CW/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
