# ham-cw.service
[Unit]
Description=ham-cw iambic CW keyer
After=local-fs.target sound.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Ham-CW
ExecStart=/usr/bin/python3 /home/pi/Ham-CW/ham_cw.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
