# ham-cw.service
[Unit]
Description=ham-cw CW keyer
After=network.target pigpiod.service
Wants=pigpiod.service

[Service]
Type=simple
WorkingDirectory=/home/pi/Ham-CW
ExecStart=/usr/bin/python3 /home/pi/Ham-CW/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
